# templates/squid-ssl.path.j2
# Ansible manageged. Local modifactions may overwritten.

[Path]
PathChanged={{ ssl_crt_path }}

[Install]
WantedBy=multi-user.target
